# phpcaptcha/supervisord.conf
[inet_http_server]
port = 127.0.0.1:9001

; Deamon
[supervisord]
logfile = /var/log/supervisord.log
logfile_maxbytes = 50MB
logfile_backups = 5
loglevel = info
pidfile = /tmp/supervisord.pid
nodaemon = false
minfds = 1024
minprocs = 200
user = root
directory = /tmp
strip_ansi = true

; RPC Interface
[rpcinterface:supervisor]
supervisor.rpcinterface_factory  =  supervisor.rpcinterface:make_main_rpcinterface

; CTL
[supervisorctl]
serverurl = http://127.0.0.1:9001

[include]
files = /etc/supervisor.d/*.ini

[program:php]
command = /usr/bin/php-fpm
user = root
startsecs = 0
autorestart = false

[program:nginx]
command = /usr/sbin/nginx
user = root
startsecs = 0
autorestart = false


